# scripts/launchd/com.cromwell.zoho-poll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cromwell.zoho-poll</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/curl</string>
    <string>-sS</string>
    <string>-o</string>
    <string>/Users/majidaljassas/cromwell-os/logs/zoho-poll.log</string>
    <string>http://localhost:3000/api/zoho/poll-bills</string>
  </array>

  <!-- Run nightly at 03:15 (after most invoices have settled in Zoho) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>15</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/majidaljassas/cromwell-os/logs/zoho-poll.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/majidaljassas/cromwell-os/logs/zoho-poll.err</string>
</dict>
</plist>
